# record-and-playback/core/systemd/avx_record_core.slice
[Unit]
Description = Averox Recording Processing
