# build/packages-template/avx-config/averox.target
[Unit]
Description=Averox System
Requires=
After=syslog.target network.target

[Install]
WantedBy=multi-user.target
